# data/ayatana-indicator-bluetooth.service.in
[Unit]
Description=Ayatana Indicator Bluetooth Service
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/ayatana-indicator-bluetooth/ayatana-indicator-bluetooth-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
