# data/indicator-bluetooth.service.in
[Unit]
Description=Indicator Bluetooth Backend
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/indicator-bluetooth-service
Restart=on-failure
